# deploy/tetracubed-fox.service
[Unit]
Description=Tetracubed Fox Discord Bot
# Ordering only — a --user unit can't hard-depend on the system network target,
# but Restart=always covers the case where it starts before the network is up.
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/tetracubed-fox
ExecStart=%h/tetracubed-fox/deploy/run.sh
Restart=always
RestartSec=10

[Install]
# default.target is the user-session equivalent of multi-user.target.
WantedBy=default.target
